# go2rtc.service
[Unit]
Description=go2rtc Camera Streaming Service
After=network.target

[Service]
Type=simple
User=USER_PLACEHOLDER
WorkingDirectory=WORKING_DIR_PLACEHOLDER
ExecStart=WORKING_DIR_PLACEHOLDER/bin/go2rtc -c WORKING_DIR_PLACEHOLDER/bin/go2rtc.yaml
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
